# config/autostart/blueman-tray.desktop
[Desktop Entry]
Exec=/bin/blueman-tray
Name=blueman-tray
Comment=bluetooth manager GUI applet
Type=Application
