# config/autostart/ibus-daemon.desktop
[Desktop Entry]
Exec=/usr/bin/ibus start
Name=ibus-daemon
Comment=daemon program for ibus
Type=Application
