# config/autostart/syncthing.desktop
[Desktop Entry]
Exec=/usr/bin/syncthing serve
Name=syncthing
Comment=Syncthing
Type=Application
