# config/autostart/thunderbird.desktop
[Desktop Entry]
Exec=/usr/bin/flatpak run org.mozilla.Thunderbird
Name=thunderbird
Comment=Thunderbird email application
Type=Application
